# conf/system/systemd/polaric-aprsd-backup.service
[Unit]
Description=Polaric APRSD daily backup
Documentation=man:polaric-aprsd(8)

[Service]
Type=oneshot
ExecStart=/etc/cron.daily/polaric-aprsd
